# stow/systemd/.config/systemd/user/album-picker-refresh.service
[Unit]
Description=Refresh album picker cache from dashboard API

[Service]
Type=oneshot
ExecStart=%h/.local/bin/album-picker-refresh
