# stow/systemd/.config/systemd/user/dashboard-cycling-1500.timer
[Unit]
Description=Dashboard cycling events at 3:00 PM

[Timer]
OnCalendar=*-*-* 15:00:00
Persistent=true
Unit=dashboard-cycling.service

[Install]
WantedBy=timers.target
