# stow/systemd/.config/systemd/user/dashboard-digest.service
[Unit]
Description=Dashboard digest curation

[Service]
Type=oneshot
WorkingDirectory=/home/jmckible/dev/dashboard
Environment=PATH=/home/jmckible/.local/bin:/home/jmckible/.local/share/mise/installs/ruby/3.4.7/bin:/usr/local/bin:/usr/bin
ExecStart=/home/jmckible/dev/dashboard/bin/agent digest
StandardOutput=append:/home/jmckible/dev/dashboard/log/agent.log
StandardError=append:/home/jmckible/dev/dashboard/log/agent.log
TimeoutStartSec=600
