# stow/systemd/.config/systemd/user/omarchy-wallpaper-auto.timer
[Unit]
Description=Omarchy Automatic Wallpaper Rotation Timer
After=graphical-session.target

[Timer]
# Run every 30 minutes
OnBootSec=1min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
