# stow/systemd/.config/systemd/user/serial-generate-0850.timer
[Unit]
Description=Serial fiction generation at 8:50 AM

[Timer]
OnCalendar=*-*-* 08:50:00
Persistent=true
Unit=serial-generate.service

[Install]
WantedBy=timers.target
