# stow/systemd/.config/systemd/user/serial-generate@.service
[Unit]
Description=Serial fiction generation for %i

[Service]
Type=oneshot
WorkingDirectory=%h/dev/serial
ExecStart=%h/.local/bin/claude -p --allowedTools Bash,Read,Write,Glob,Grep,Agent "Run /serial:generate %i"
StandardOutput=append:%h/dev/serial/log/generate.log
StandardError=append:%h/dev/serial/log/generate.log
