# stow/systemd/.config/systemd/user/serial-generate-1050.timer
[Unit]
Description=Serial fiction generation at 10:50 AM

[Timer]
OnCalendar=*-*-* 10:50:00
Persistent=true
Unit=serial-generate.service

[Install]
WantedBy=timers.target
