# stow/systemd/.config/systemd/user/serial-generate-1450.timer
[Unit]
Description=Serial fiction generation at 2:50 PM

[Timer]
OnCalendar=*-*-* 14:50:00
Persistent=true
Unit=serial-generate@the-broken-oath.service

[Install]
WantedBy=timers.target
